# mydependent.service
[Unit]
Description=Dependent Service
After=myfirst.service
Requires=myfirst.service

[Service]
Type=simple
ExecStart=/bin/echo "Dependent service started"
ExecStop=/bin/echo "Dependent service stopped"
SuccessAction=none

