# myenv.service
[Service]
Type=simple
ExecStart=/usr/bin/env bash -c 'echo $MY_ENV_VAR; sleep infinity'
Environment=MY_ENV_VAR=HelloWorld
ExecStop=/bin/echo "Stopping myenv.service"
SuccessAction=none
