# myrestart.service
[Unit]
Description=My Restart Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'sleep infinity'
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
