# mytimer.service
[Unit]
Description=My Timer Service

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "Timer fired!"

